# sqlpage.service
# SQLPage systemd service file for manual/source installations
# For package installations, see debian/sqlpage.service
# Install to: /etc/systemd/system/sqlpage.service

[Unit]
Description=SQLPage Web Server
Documentation=https://sql-page.com
After=network.target

[Service]
Type=simple
User=sqlpage
Group=sqlpage
WorkingDirectory=/var/www/sqlpage
ExecStart=/usr/local/bin/sqlpage.bin

Environment="RUST_LOG=info"
Environment="LISTEN_ON=0.0.0.0:80"

AmbientCapabilities=CAP_NET_BIND_SERVICE

Restart=always
RestartSec=10s

SyslogIdentifier=sqlpage

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/sqlpage /var/log/sqlpage
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
